# etc/xinetd.d/alternc_lxc
service alternc_lxc
{
	type           = UNLISTED
	port           = 6504
	socket_type    = stream
	protocol       = tcp
	wait           = no
	user           = root
	server         = /var/lib/alternc-lxc/alternc_lxc_remote

	only_from      = 127.0.0.1 

	disable        = no
}


