# debian/alternc-awstats.cron.d
30 4 * * *      root    /usr/lib/alternc/alternc-awstats
